# esnode/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=/startNode.sh
autorestart=false
priority=30

[program:install-plugin]
command=/plugin.sh
autorestart=false
priority=50



